# celery.conf
[program:celery]

command=/home/mary/PycharmProjects/university/env/bin/celery -A students_tracker worker -l info
directory=/home/mary/PycharmProjects/university/src
autostart=true
autorestart=true

