# uwshi.conf
[program:uwsgi]

command=/home/mary/PycharmProjects/university/env/bin/uwsgi uwsgi.ini
directory=/home/mary/PycharmProjects/university/
autostart=true
autorestart=true